# config/live-build/config/includes.chroot/etc/systemd/system/speedos-first-boot.service
[Unit]
Description=SpeedOS First Boot Setup
After=network.target lightdm.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/speedos-first-boot.sh
StandardOutput=journal

[Install]
WantedBy=multi-user.target
