# config/live-build/config/includes.chroot/etc/systemd/system/speedos-update.timer
[Unit]
Description=Run SpeedOS Update Daily

[Timer]
OnBootSec=15min
OnUnitActiveSec=24h

[Install]
WantedBy=timers.target
