# config/live-build/config/includes.chroot/etc/systemd/system/speed-ai-monitor.service
[Unit]
Description=SpeedOS AI Monitor and Optimizer
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/speedos/speedboost/speed-ai-monitor.py
Restart=always
User=root

[Install]
WantedBy=multi-user.target
