# airootfs/etc/xdg/autostart/speedos-firstboot.desktop
[Desktop Entry]
Type=Application
Exec=xfce4-firstboot
Hidden=false
